# docker/app/root/etc/s6-overlay/s6-rc.d/init-security-check/up
/etc/s6-overlay/scripts/init-security-check
