# docker/app/root/etc/s6-overlay/s6-rc.d/svc-php-fpm/run
#!/command/execlineb -P
with-contenv
/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
