# docker/app/root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
with-contenv
nginx -g "daemon off;"