# docker/app/root/etc/s6-overlay/s6-rc.d/svc-scheduler/run
#!/command/execlineb -P

with-contenv
php /var/www/html/bin/console messenger:consume scheduler_default --time-limit=300 --memory-limit=128M --limit=10 --env=`printcontenv APP_ENV` --quiet
