# docker/app/root/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
with-contenv
/usr/local/sbin/php-fpm -F -O